# crontab
30 3 */7 * * /app/update_db.sh > /proc/1/fd/1 2>/proc/1/fd/2
